# config/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm
autostart=true
autorestart=true
stderr_logfile=/var/log/php-fpm.err.log
stdout_logfile=/var/log/php-fpm.out.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log

[program:streamlit]
command=/usr/local/bin/streamlit run /ai/chat/app.py
directory=/ai/chat
autostart=true
autorestart=true
stderr_logfile=/var/log/streamlit.err.log
stdout_logfile=/var/log/streamlit.out.log

[program:watcher]
command=python3 /ai/chat/watcher.py
directory=/ai/chat
autostart=true
autorestart=true
stderr_logfile=/var/log/watcher.err.log
stdout_logfile=/var/log/watcher.out.log

[program:run-script]
command=/bin/bash /ai/chat/run.sh
directory=/ai/chat
autostart=true
autorestart=true
stderr_logfile=/var/log/run-script.err.log
stdout_logfile=/var/log/run-script.out.log

[include]
files = /etc/supervisor/conf.d/*.conf
